# aidialogs-watcher.service
[Unit]
Description=AI Dialogs Git Watcher
After=network.target

[Service]
Type=simple
User=semion
WorkingDirectory=/opt/aidialogs
ExecStart=/opt/aidialogs/scripts/git-watcher.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

