# aidialogs-bot.service
[Unit]
Description=AI Dialogs Telegram Bot
After=network.target

[Service]
Type=simple
User=semion
WorkingDirectory=/opt/aidialogs
ExecStart=/home/semion/.local/bin/uv run python -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

